# linux/data/cmcc@.service.in
[Unit]
Description=Networking for netctl profile cmcc on %I
Documentation=man:netctl.profile(5)
Before=network.target netctl.service
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=cmugdir@cmug -i %I
ExecStartPre=ip@ link set %I down
ExecStartPre=cp@ /usr/share/cmug/pppoe.conf /etc/netctl/cmcc
ExecStart=/usr/lib/network/network start cmcc
ExecStop=/usr/lib/network/network stop cmcc
